# src/uni-drop/tests/e2e/supervisord.conf
[supervisord]
nodaemon=true
user=root


[program:Xvfb]
command=/usr/bin/Xvfb %(ENV_DISPLAY)s -screen 0 %(ENV_GEOMETRY)s
# command=Xvfb %(ENV_DISPLAY)s -screen 0 $GEOMETRY -ac +extension RANDR
autorestart=true
startretries=1
stdout_logfile=/workspace/xvfb.stdout.log
stderr_logfile=/workspace/xvfb.stderr.log

[program:x11vnc]
command=/usr/bin/x11vnc -display %(ENV_DISPLAY)s -nopw -listen localhost -xkb -ncache 10 -ncache_cr -forever
stdout_logfile=/workspace/x11vnc.stdout.log
stderr_logfile=/workspace/x11vnc.stderr.log
autorestart=true
depends_on=Xvfb
startretries=1

[program:noVNC]
command=/workspace/noVNC/utils/novnc_proxy --vnc localhost:5900 --listen localhost:6080
autorestart=true
startretries=1
stdout_logfile=/workspace/noVNC.stdout.log
stderr_logfile=/workspace/noVNC.stderr.log

[program:fluxbox]
command=/usr/bin/fluxbox -display %(ENV_DISPLAY)s > /workspace/fluxbox.stdout.log
autorestart=true
startretries=1
stdout_logfile=/workspace/fluxbox.stdout.log
stderr_logfile=/workspace/fluxbox.stderr.log